# deploy/gunicorn.service
[Unit]
Description=Contact360 DocsAI Admin (Gunicorn) — simple HTTP bind template (127.0.0.1:8000)
After=network.target
# Production: use deploy/systemd/gunicorn.service + gunicorn.socket (Unix socket to Nginx).
# WSGI entrypoint matches deploy/systemd/gunicorn-start.sh: docsai.wsgi:application, DJANGO_SETTINGS_MODULE=config.settings.

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/docsai
Environment="DJANGO_SETTINGS_MODULE=config.settings"
EnvironmentFile=/home/ubuntu/docsai/.env.prod
ExecStart=/home/ubuntu/docsai/venv/bin/gunicorn docsai.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --timeout 120 \
    --access-logfile /var/log/docsai/access.log \
    --error-logfile /var/log/docsai/error.log
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
