# deploy/gunicorn.service
[Unit]
Description=Contact360 DocsAI Admin (Gunicorn) — simple HTTP bind template
After=network.target
# For production with Nginx + Unix socket, use deploy/systemd/gunicorn.service instead.

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/docsai
Environment="DJANGO_SETTINGS_MODULE=docsai.settings"
EnvironmentFile=/home/ubuntu/docsai/.env.prod
ExecStart=/home/ubuntu/docsai/venv/bin/gunicorn docsai.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --timeout 120 \
    --access-logfile /var/log/docsai/access.log \
    --error-logfile /var/log/docsai/error.log
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
